# backups/20260225_142108_v5_allinone/trident.service
[Unit]
Description=Trident Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trident-bot
EnvironmentFile=/home/ubuntu/trident-bot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/trident-bot/venv/bin/python -u /home/ubuntu/trident-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
